# runtime/supervisord.conf
[supervisord]
nodaemon=true

[program:runtime]
command=java -jar /agent-runtime/app.jar